# system-services/aether-recovery.service
[Unit]
Description=AetherOS Recovery Agent
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/aetheros/recovery/recovery_agent.py

[Install]
WantedBy=recovery.target
